# dot_config/systemd/user/dunst-update.path
[Path]
PathModified=%h/.config/dunst/dunstrc

[Install]
WantedBy=default.target
