# dot_config/systemd/user/polybar-update.path
[Path]
PathModified=%h/.config/polybar/config.ini

[Install]
WantedBy=default.target
