# dot_config/systemd/user/rss.service
[Unit]
Description="Newsboat notifications"
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -x reload
