# dot_config/polybar/polybar-update.path
[Path]
PathModified=%h/.config/polybar/config.ini

